# watchdog.service
[Unit]
Description=Elon Bot Watchdog
After=elon-bot.service
Wants=elon-bot.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/elon
ExecStart=/root/elon/.venv/bin/python /root/elon/bot_watchdog.py --interval 60
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target