# watchdog.service
[Unit]
Description=Elon Bot Watchdog Service
After=elon-bot.service
Wants=elon-bot.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/elon
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/root/elon
ExecStart=/usr/bin/python3 /root/elon/bot_watchdog.py --interval 60 --max-restarts 3
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=elon-watchdog

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/elon

# Resource limits
LimitNOFILE=65536
MemoryMax=128M

[Install]
WantedBy=multi-user.target
